# launchd/com.franlledo.seo-refresh-recommender.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.franlledo.seo-refresh-recommender</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>cd "/Users/franlledo/Documents/Claude/Projects/seo-franlledo" &amp;&amp; ./coordinator.sh refresh-recommender</string>
    </array>

    <!-- Primer día de cada mes a las 7:00. Phase 1 (Python) escanea sitemap +
         consulta GSC URL Inspection para cada post. Phase 2 (Claude) lee los
         flags y produce refresh-queue.json con la acción recomendada por URL
         (refresh, expand, redirect, request_indexing, wait). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/seo-refresh-recommender.stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/seo-refresh-recommender.stderr</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
